# mydns/mydns.timer
[Unit]
Description=mydns notify timer

[Timer]
OnCalendar=*-*-* 6:00
Unit=mydns.service

[Install]
WantedBy=timers.target
