# util/system/mydns.timer
[Unit]
Description=mydns notice service

[Timer]
OnCalendar=*-*-* 5:10:00
Unit=mydns.service

[Install]
WantedBy=timers.target