# util/system/mydns.timer
[Unit]
Description=mydns notice service

[Timer]
OnCalendar=*-*-* 4:00:00
Unit=mydns.service

[Install]
WantedBy=timers.target