# util/system/reboot.timer
[Unit]
Description=minecraft reboot timer

[Timer]
OnCalendar=*-*-* 6:00
Unit=reboot.service

[Install]
WantedBy=timers.target